# xsstatus.cron
#!/bin/sh
/usr/local/sbin/xsstatus >/etc/motd
